# periodic/daily/900.sanoid
TZ=UTC
lockf /var/run/sanoid.lock sanoid --take-snapshots --verbose
lockf /var/run/sanoid.lock sanoid --prune-snapshots --verbose
