# periodic/hourly/900.sanoid
TZ=UTC
lockf /var/run/sanoid.lock sanoid --take-snapshots --quiet
lockf /var/run/sanoid.lock sanoid --prune-snapshots --quiet
